# packaging/systemd/rsyncd.service
[Unit]
Description=rsync daemon providing rsync protocol services
Documentation=https://github.com/oferchen/rsync
After=network.target

[Service]
Type=simple
Environment="RSYNCD_CONFIG=/etc/rsyncd.conf"
EnvironmentFile=-/etc/default/rsyncd
ExecStart=/usr/bin/rsyncd --config ${RSYNCD_CONFIG} $RSYNCD_ARGS
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
TimeoutStopSec=15s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=CAP_CHOWN CAP_DAC_OVERRIDE CAP_FOWNER CAP_FSETID CAP_MKNOD CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE
LimitNOFILE=16384
UMask=0077

[Install]
WantedBy=multi-user.target
